# bot.conf
[program:telegrambot]
directory=/bots/parser_bot
command=python3 app.py

autostart=true
autorestart=true
std_err_logfile=/bots/parser_bot/logfile_err.logfile
